# ansible/roles/docker_compose/files/docker-resource-limit.slice
[Unit]
Description=Slice that limits docker resources
Before=slices.target

[Slice]
CPUAccounting=true
# Allow usage of 1.9 cpu cores
CPUQuota=190%

# Memory Management
MemoryAccounting=true
# Desired maximum ram usage
MemoryHigh=1536M
# Maximum ram usage
MemoryMax=1664M
# 1900M at all => 236M left for system
